# xroad-catalog-lister/packages/xroad-catalog-lister/debian/xroad-catalog-lister.postinst
#!/bin/sh

if ! id xroad-catalog > /dev/null 2>&1 ; then

    adduser --system --no-create-home \
                --group --disabled-password --shell /bin/false \
                xroad-catalog
fi                
chmod 755 /usr/lib/xroad-catalog/xroad-catalog-lister.jar
chown -R xroad-catalog:xroad-catalog /usr/lib/xroad-catalog


invoke-rc.d xroad-catalog-lister start