# xroad-catalog-collector/packages/xroad-catalog-collector/redhat/SOURCES/xroad-catalog-collector.cron
SHELL=/bin/bash
0 1 * * * root xroad-catalog service xroad-catalog-collector restart


