# xroad-catalog-collector/packages/xroad-catalog-collector/redhat/SOURCES/xroad-catalog-collector.cron
#!/bin/sh
/usr/share/xroad/bin/xroad-catalog-collector

