# xroad-catalog-collector/packages/xroad-catalog-collector/redhat/SOURCES/xroad-catalog-collector.service
[Unit]
Description=X-Road Catalog Collector
Wants=postgresql-9.6.service
After=network.target postgresql-9.6.service

[Service]
PermissionsStartOnly=true
User=xroad
Group=xroad
ExecStart=/usr/share/xroad/bin/xroad-catalog-collector
EnvironmentFile=/etc/xroad/catalog-profile.properties
Restart=on-failure

[Install]
WantedBy=multi-user.target

