# xroad-catalog-collector/packages/xroad-catalog-collector/redhat/SOURCES/xroad-catalog-collector.service
[Unit]
Description=X-Road Catalog Collector
Wants=postgresql*.service
After=network.target postgresql*.service

[Service]
PermissionsStartOnly=true
User=xroad-catalog
Group=xroad-catalog
ExecStart=/usr/share/xroad/bin/xroad-catalog-collector
Restart=on-failure

[Install]
WantedBy=multi-user.target

